# packaging/debian/prerm
#!/bin/bash
# NV Broadcast pre-removal script
set -e

case "$1" in
    upgrade|deconfigure)
        # On package upgrade, stop the running app but preserve service state.
        pkill -f "nvbroadcast" 2>/dev/null || true
        if systemctl --user is-active nvbroadcast-vcam.service &>/dev/null; then
            systemctl --user stop nvbroadcast-vcam.service 2>/dev/null || true
        fi
        if systemctl --user is-active nvbroadcast-audio.service &>/dev/null; then
            systemctl --user stop nvbroadcast-audio.service 2>/dev/null || true
        fi
        ;;
    remove|purge)
        pkill -f "nvbroadcast" 2>/dev/null || true
        if systemctl --user is-active nvbroadcast-vcam.service &>/dev/null; then
            systemctl --user stop nvbroadcast-vcam.service 2>/dev/null || true
        fi
        if systemctl --user is-active nvbroadcast-audio.service &>/dev/null; then
            systemctl --user stop nvbroadcast-audio.service 2>/dev/null || true
        fi
        if systemctl --user is-enabled nvbroadcast-vcam.service &>/dev/null; then
            systemctl --user disable nvbroadcast-vcam.service 2>/dev/null || true
        fi
        if systemctl --user is-enabled nvbroadcast-audio.service &>/dev/null; then
            systemctl --user disable nvbroadcast-audio.service 2>/dev/null || true
        fi
        ;;
esac

#DEBHELPER#
